# qr-code-bot.service
[Unit]
Description=QR Code Generator Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/qr_code_bot
Environment="PATH=/opt/qr_code_bot/venv/bin"
ExecStart=/opt/qr_code_bot/venv/bin/python /opt/qr_code_bot/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Ограничения ресурсов
LimitNOFILE=65536
MemoryMax=512M

# Безопасность
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

